# recipes-debug/vendor-tools/files/network_setup.service
[Unit]
Description=Network Setup Service
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/network_setup.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
